# services/docker-logger.service
[Unit]
Description=Docker LRU Tracker
After=network.target

[Service]
Type=simple
ExecStart=/root/.docker-history/scripts/logger.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

